# docker/etc/s6/confluence/run
#!/bin/bash

# 配置
source ./setup

# 启动Jira服务器
JAVA_OPTS="${JAVA_OPTS} -javaagent:/opt/atlassian/agent/agent.jar"
exec gosu "${USERNAME}" /opt/atlassian/confluence/bin/start-confluence.sh -fg
